# install/supervisor.conf
[supervisord]

nodaemon=true

[program:nginx]

user=root
command=/usr/sbin/nginx -g "daemon off;"
autorestart=false
